# deploy/bare-metal/systemd/thinkq-notifications.service
[Unit]
Description=ThinkQ Notifications Service
After=network.target redis-server.service thinkq-tickets.service
Wants=redis-server.service thinkq-tickets.service

[Service]
Type=simple
WorkingDirectory=/opt/thinkq/work_backend/notifications-service
EnvironmentFile=/opt/thinkq/env/notifications-service.env
ExecStart=/usr/bin/node src/server.js
Restart=always
RestartSec=5
User=www-data
Group=www-data

[Install]
WantedBy=multi-user.target
