# THINKQ-baremetal/deploy/bare-metal/systemd/thinkq-data.service
[Unit]
Description=ThinkQ Data Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/thinkq/backend/data-service
EnvironmentFile=/opt/thinkq/env/data-service.env
Environment=JAVA_TOOL_OPTIONS=-Dfile.encoding=UTF-8
Environment=JAVA_HOME=/usr/lib/jvm/java-21-amazon-corretto
ExecStart=/usr/lib/jvm/java-21-amazon-corretto/bin/java -jar /opt/thinkq/backend/data-service/target/data-service-1.0.0.jar
Restart=always
RestartSec=5
User=thinkq
Group=thinkq
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target